# systemd/raspi-openclaw-ops.service
[Unit]
Description=raspi-openclaw-ops status server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/raspi-openclaw-ops
ExecStart=/usr/bin/node /opt/raspi-openclaw-ops/dist/server.js
Restart=always
RestartSec=2
Environment=PORT=8080
Environment=HOST=0.0.0.0
# Health check (optional)
# Use ONE of the following:
# 1) systemd unit check
# Environment=CLAWDBOT_SERVICE=clawdbot-gateway
# 2) process existence check (pgrep -f)
# Environment=CLAWDBOT_PROCESS_PATTERN=clawdbot-gateway

Environment=CLAWDBOT_SERVICE=
Environment=CLAWDBOT_PROCESS_PATTERN=

# Hardening (optional)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/raspi-openclaw-ops

[Install]
WantedBy=multi-user.target
